# roles/server_metrics/templates/my-service.service.j2
[Unit]
Description=Server metrics collector

[Service]
Type=oneshot
ExecStart=/usr/local/bin/metrics.sh
